# luci-mod-39wrt-luci-custom/files/etc/uci-defaults/99-39wrt-luci-custom
#!/bin/sh

[ -x /usr/libexec/39wrt-luci-custom/apply ] || exit 0
/usr/libexec/39wrt-luci-custom/apply apply
exit $?
